# api.service
[Unit]
Description=API server

[Service]
ExecStart=/tmp/api

[Install]
WantedBy=multi-user.target
Alias=api.service
